# net-core-3-services/testapp.service
[Unit]
Description=test app

[Service]
Type=notify
ExecStart=/usr/sbin/testapp/NetCore3Service

[Install]
WantedBy=multi-user.target
